# scripts/cloudflared.service
[Unit]
Description=Cloudflare Tunnel (TaskHive API)
After=network.target taskhive-api.service
Wants=taskhive-api.service

[Service]
User=root
ExecStart=/usr/bin/cloudflared tunnel --url http://localhost:8000 --logfile /var/log/cloudflared.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
