# scripts/taskhive-api.service
[Unit]
Description=TaskHive API
After=network.target

[Service]
User=root
WorkingDirectory=/opt/taskhive/repo
EnvironmentFile=/opt/taskhive/repo/.env
ExecStart=/opt/taskhive/repo/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
