# scripts/taskhive-reviewer.service
[Unit]
Description=TaskHive Reviewer Agent
After=network.target taskhive-api.service
Wants=taskhive-api.service

[Service]
User=root
WorkingDirectory=/opt/taskhive/repo
EnvironmentFile=/opt/taskhive/repo/.env
Environment=PYTHONUNBUFFERED=1
Environment=TASKHIVE_BASE_URL=http://127.0.0.1:8000
# TASKHIVE_REVIEWER_API_KEY is loaded from EnvironmentFile (.env) directly — do NOT use ${VAR} here
ExecStart=/opt/taskhive/repo/.venv/bin/python -u reviewer-agent/run.py \
    --daemon \
    --interval 30
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
